# supervisor/idle.conf
[program:idl]
command=/tmp/idle.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/idle.err.log
stdout_logfile=/var/log/idle.out.log
